# scripts/community-view.service
[Unit]
Description=Community View Backend Manager
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/comunity_view_backend
Environment=PATH=/home/ubuntu/comunity_view_backend/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/comunity_view_backend/venv/bin/python /home/ubuntu/comunity_view_backend/community_view_manager.py daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=community-view

[Install]
WantedBy=multi-user.target 